# aports/main/fixtime/fixtime.pre-deinstall
#!/bin/sh

rc-update -q del fixtime default
rc-update -q del fixtime shutdown
