# aports/main/fixtime/fixtime.post-install
#!/bin/sh

rc-update -q add fixtime default
rc-update -q add fixtime shutdown
