# aports/main/fixtime/fixtime.post-install
#!/bin/sh

rc-update -q add fixtime boot
rc-update -q add fixtime shutdown